# ubuntu-lxde/supervisord.conf
[supervisord]
nodaemon=false

[program:sshd]
command=/usr/bin/sshd -D
